# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[program:uvicorn]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8002
directory=/app/backend
stdout_logfile=/var/log/supervisor/uvicorn.log
stderr_logfile=/var/log/supervisor/uvicorn_error.log
autorestart=true
startsecs=10
priority=5

[program:nginx]
command=/usr/sbin/nginx -c /app/nginx.conf
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log
autorestart=true
startsecs=5
priority=10
